# ansible/roles/supervisord/templates/supervisord_config.conf
[program:{{item.name}}]
command={{item.command}}
autostart={{item.autostart}}
autorestart={{item.autorestart}}
